# examples/example2/https.socket
[Socket]
ListenDatagram=443
ListenStream=443
FileDescriptorName=websecure
Service=traefik.service

[Install]
WantedBy=sockets.target
