# examples/example1/https.socket
[Socket]

ListenStream=0.0.0.0:443
FileDescriptorName=websecure
Service=mytraefik.service

[Install]
WantedBy=sockets.target
