# examples/example1/http.socket
[Socket]

BindIPv6Only=both
FileDescriptorName=web
ListenStream=[::]:80
Service=mytraefik.service

[Install]
WantedBy=sockets.target
