# examples/example1/https.socket
[Socket]

BindIPv6Only=both
FileDescriptorName=websecure
ListenStream=[::]:443
Service=mytraefik.service

[Install]
WantedBy=sockets.target
