# examples/example1/https.socket
[Socket]

ListenStream=443
FileDescriptorName=websecure
Service=traefik.service

[Install]
WantedBy=sockets.target
