# examples/example1/http.socket
[Socket]

ListenStream=80
FileDescriptorName=web
Service=mytraefik.service

[Install]
WantedBy=sockets.target
