# examples/example1/mynginx.container
[Container]
Image=docker.io/library/nginx
ContainerName=mynginx
Network=mynet.network
Label="traefik.enable=true"
