# examples/example1/whoami.container
[Container]
Image=docker.io/traefik/whoami
ContainerName=whoami.example.com
Network=mynet.network
Label="traefik.enable=true"
